# ai-job-processor.service
[Unit]
Description=Portavi AI Job Processor
After=network.target

[Service]
User=root
WorkingDirectory=/var/www/portavi
ExecStart=/usr/bin/npm run process-ai-jobs
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ai-job-processor

[Install]
WantedBy=multi-user.target
